# sstat/src/sstat-post.cron
#!/bin/bash

PATH=/bin:/usr/bin:/sbin:/usr/sbin

sstat-post
