# lcopy/trunk/etc/lcopy-refresh.cron
#!/bin/sh

for x in /etc/lcopy-refresh.d/*; do
    if test -f "${x}"; then
	LCOPY_USER=root
	LCOPY_LOG=/dev/null
	LCOPY_DIR=.
	source "${x}"

	if test ! -f "${LCOPY_LOG}"; then
	    LCOPY_LOG=/dev/null
	fi

	if test -d ${LCOPY_DIR}; then
	    su  \
		-c "lcopy-refresh --log=${LCOPY_LOG} ${LCOPY_DIR}" \
		- ${LCOPY_USER} > /dev/null 2>&1
	fi
    fi
done


